# container/supervisor/atlas.ini
[program:atlas]
command=/usr/bin/start_atlas.sh
autostart=true
autorestart=false
