# container/supervisor/zookeeper.ini
[program:zookeeper]
command=/opt/atlas/zk/bin/zkServer.sh start-foreground
autostart=true
autorestart=false
priority=10
